# void-stuff/common/etc/sv/nix-daemon/run
#!/bin/sh
exec 2>&1
exec /root/.nix-profile/bin/nix-daemon --option extra-trusted-users @wheel
